# deploy/systemd/zunder-zapfe-web.service.in
[Unit]
Description=Zunder Zapfe local web application
After=local-fs.target

[Service]
Type=simple
User=@@SERVICE_USER@@
Group=@@SERVICE_USER@@
WorkingDirectory=@@APP_DIR@@
Environment=ZUNDER_ZAPFE_DATABASE_URL=sqlite:////var/lib/zunder-zapfe/zunder-zapfe.db
EnvironmentFile=-/etc/zunder-zapfe/web.env
ExecStartPre=@@APP_DIR@@/.venv/bin/alembic -c @@APP_DIR@@/alembic.ini upgrade head
ExecStart=@@APP_DIR@@/.venv/bin/zunder-zapfe
Restart=on-failure
RestartSec=2
TimeoutStopSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/zunder-zapfe /var/log/zunder-zapfe

[Install]
WantedBy=multi-user.target
